# deploy/docker/supervisord.conf
[supervisord]
user=breaker
logfile=%(ENV_CB_DATA_DIR)s/supervisord.log
pidfile=%(ENV_CB_DATA_DIR)s/supervisord.pid
nodaemon=true
childlogdir=%(ENV_CB_DATA_DIR)s

; ------------------------------------------------------------------------------
; Postgres (embedded; Debian Bookworm default is PostgreSQL 15)
; Wrapper: only start postgres when CB_DB_URL is localhost/127.0.0.1; otherwise sleep (external DB).
; ------------------------------------------------------------------------------
[program:postgres]
priority=10
command=/opt/circuitbreaker/deploy/common/run-postgres-if-embedded.sh
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=30
stderr_logfile=%(ENV_CB_DATA_DIR)s/pg_err.log
stdout_logfile=%(ENV_CB_DATA_DIR)s/pg.log
user=breaker
redirect_stderr=false

; ------------------------------------------------------------------------------
; pgbouncer — connection pooler (transaction-mode pooling for the backend)
; Sits between backend and Postgres on port 6432.
; Only starts when CB_DB_URL routes through localhost (embedded DB mode).
; ------------------------------------------------------------------------------
[program:pgbouncer]
priority=15
command=/bin/sh -c 'if echo "$CB_DB_URL" | grep -qE "localhost|127\\.0\\.0\\.1"; then exec /usr/sbin/pgbouncer /etc/pgbouncer/pgbouncer.ini; else echo "[pgbouncer] External DB — sleeping"; exec sleep infinity; fi'
directory=%(ENV_CB_DATA_DIR)s
autostart=true
autorestart=true
startsecs=2
startretries=10
stderr_logfile=%(ENV_CB_DATA_DIR)s/pgbouncer_err.log
stdout_logfile=%(ENV_CB_DATA_DIR)s/pgbouncer.log
user=breaker

; ------------------------------------------------------------------------------
; NATS JetStream (embedded; auth mandatory via NATS_AUTH_TOKEN)
; ------------------------------------------------------------------------------
[program:nats]
priority=20
command=/bin/sh -c 'if [ -z "$NATS_AUTH_TOKEN" ]; then echo "[nats] FATAL: NATS_AUTH_TOKEN is required" >&2; exit 1; fi; exec nats-server --jetstream --store_dir %(ENV_CB_DATA_DIR)s/nats --auth "$NATS_AUTH_TOKEN"'
directory=%(ENV_CB_DATA_DIR)s
autostart=true
autorestart=true
startsecs=3
startretries=10
stderr_logfile=%(ENV_CB_DATA_DIR)s/nats_err.log
stdout_logfile=%(ENV_CB_DATA_DIR)s/nats.log
user=breaker

; ------------------------------------------------------------------------------
; Redis (embedded; telemetry cache + pub/sub for real-time push)
; Memory capped at 128 MB with LRU eviction for cache safety.
; ------------------------------------------------------------------------------
[program:redis]
priority=25
command=/bin/sh -c 'PASS=$(cat %(ENV_CB_DATA_DIR)s/.redis_pass 2>/dev/null || echo ""); if [ -n "$PASS" ]; then exec /usr/bin/redis-server --daemonize no --dir %(ENV_CB_DATA_DIR)s/redis --maxmemory 128mb --maxmemory-policy allkeys-lru --bind 127.0.0.1 --protected-mode yes --save "" --requirepass "$PASS"; else exec /usr/bin/redis-server --daemonize no --dir %(ENV_CB_DATA_DIR)s/redis --maxmemory 128mb --maxmemory-policy allkeys-lru --bind 127.0.0.1 --protected-mode yes --save ""; fi'
directory=%(ENV_CB_DATA_DIR)s
autostart=true
autorestart=true
startsecs=2
startretries=10
stderr_logfile=%(ENV_CB_DATA_DIR)s/redis_err.log
stdout_logfile=%(ENV_CB_DATA_DIR)s/redis.log
user=breaker

; ------------------------------------------------------------------------------
; FastAPI backend API (Uvicorn)
; Inherits CB_DB_URL, CB_DATA_DIR, NATS_URL, CB_VAULT_KEY, etc. from container env.
; ------------------------------------------------------------------------------
[program:backend-api]
priority=30
command=/usr/local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2 --no-access-log
directory=%(ENV_CB_APP_ROOT)s/backend
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=%(ENV_CB_DATA_DIR)s/backend_api.log
stderr_logfile=%(ENV_CB_DATA_DIR)s/backend_api_err.log
user=breaker

; ------------------------------------------------------------------------------
; Workers — named programs instead of numprocs for clarity
; ------------------------------------------------------------------------------
[program:worker-discovery]
priority=30
command=/usr/bin/python -m app.workers.main --type discovery
directory=%(ENV_CB_APP_ROOT)s/backend
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=%(ENV_CB_DATA_DIR)s/worker_discovery.log
stderr_logfile=%(ENV_CB_DATA_DIR)s/worker_discovery_err.log
user=breaker

[program:worker-webhook]
priority=30
command=/usr/bin/python -m app.workers.main --type webhook
directory=%(ENV_CB_APP_ROOT)s/backend
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=%(ENV_CB_DATA_DIR)s/worker_webhook.log
stderr_logfile=%(ENV_CB_DATA_DIR)s/worker_webhook_err.log
user=breaker

[program:worker-notification]
priority=30
command=/usr/bin/python -m app.workers.main --type notification
directory=%(ENV_CB_APP_ROOT)s/backend
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=%(ENV_CB_DATA_DIR)s/worker_notification.log
stderr_logfile=%(ENV_CB_DATA_DIR)s/worker_notification_err.log
user=breaker

[program:worker-telemetry]
priority=30
command=/usr/bin/python -m app.workers.main --type telemetry
directory=%(ENV_CB_APP_ROOT)s/backend
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=%(ENV_CB_DATA_DIR)s/worker_telemetry.log
stderr_logfile=%(ENV_CB_DATA_DIR)s/worker_telemetry_err.log
user=breaker

; ------------------------------------------------------------------------------
; Nginx (frontend SPA + API proxy + WebSockets). Runs as breaker on 8080/8443.
; ------------------------------------------------------------------------------
[program:nginx]
priority=40
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
startsecs=3
stdout_logfile=%(ENV_CB_DATA_DIR)s/nginx.log
stderr_logfile=%(ENV_CB_DATA_DIR)s/nginx_err.log
user=breaker
