# deploy/systemd/circuitbreaker.slice
[Unit]
Description=Circuit Breaker Service Slice

[Slice]
MemoryMax=3G
MemoryHigh=2G
