# docker/supervisord.mono.conf
[supervisord]
logfile=/data/supervisord.log
pidfile=/data/supervisord.pid
nodaemon=true
childlogdir=/data

; ------------------------------------------------------------------------------
; Postgres (embedded; Debian Bookworm default is PostgreSQL 15)
; Data dir and config live under /data so they persist across container restarts.
; ------------------------------------------------------------------------------
[program:postgres]
command=/usr/lib/postgresql/15/bin/postgres -D /data/pgdata -p 5432 -c config_file=/data/postgresql.conf
autostart=true
autorestart=true
startsecs=5
startretries=10
stderr_logfile=/data/pg_err.log
stdout_logfile=/data/pg.log
user=breaker
redirect_stderr=false

; ------------------------------------------------------------------------------
; NATS JetStream (embedded)
; Uses a small shell wrapper so we can add --auth only when NATS_AUTH_TOKEN is set.
; ------------------------------------------------------------------------------
[program:nats]
command=/bin/sh -c 'if [ -n "$NATS_AUTH_TOKEN" ]; then exec nats-server --jetstream --store_dir /data/nats --auth "$NATS_AUTH_TOKEN"; else exec nats-server --jetstream --store_dir /data/nats; fi'
directory=/data
autostart=true
autorestart=true
startsecs=3
startretries=10
stderr_logfile=/data/nats_err.log
stdout_logfile=/data/nats.log
user=breaker

; ------------------------------------------------------------------------------
; FastAPI backend API (Uvicorn)
; Inherits CB_DB_URL, CB_DATA_DIR, NATS_URL, CB_VAULT_KEY, etc. from container env.
; ------------------------------------------------------------------------------
[program:backend-api]
command=/usr/local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2 --no-access-log
directory=/app/backend
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=/data/backend_api.log
stderr_logfile=/data/backend_api_err.log
user=breaker

; ------------------------------------------------------------------------------
; Unified workers (discovery, webhook, notification)
; process_num:
;   0 → discovery
;   1 → webhook
;   2 → notification
; (Can scale further by increasing numprocs.)
; ------------------------------------------------------------------------------
[program:workers]
numprocs=3
process_name=worker-%(process_num)02d
command=/usr/bin/python -m app.workers.main --type=%(process_num)d
directory=/app/backend
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=20
stdout_logfile=/data/worker_%(process_num)d.log
stderr_logfile=/data/worker_%(process_num)d_err.log
user=breaker

; ------------------------------------------------------------------------------
; Nginx (frontend SPA + API proxy + WebSockets)
; ------------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/data/nginx.log
stderr_logfile=/data/nginx_err.log
user=breaker

